# overlay/etc/cron.daily/innobox-backup
#!/bin/sh

#run the backup script as the Apache user
su -c innobox_perform_backup www-data
